# user-lib-systemd-system/timeshift-boot.timer
[Unit]
Description="Timeshift Snapshoot on Boot"

[Timer]
OnBootSec=2min
Unit=timeshift-boot.service

[Install]
WantedBy=multi-user.target
