# systemd/system/change-machine-id.timer
[Unit]
Description="Changes machine-id on boot"

[Timer]
OnBootSec=3min
Unit=change-machine-id.service

[Install]
WantedBy=multi-user.target
