# systemd/system/change-machine-id.service
 
[Unit]
Description="A script to change machine-id"

[Service]
ExecStart=/usr/local/bin/change-machine-id.sh
